# tunnels/rootfs/etc/services.d/tunnels/finish
#!/usr/bin/env bashio

bashio::log.warning "SSH tunnel service finished. The service will be restarted automatically."